# systemd/forge-watcher.service
[Unit]
Description=Jarvis Forge queue watcher
After=network-online.target jarvis-router.service
Wants=jarvis-router.service

[Service]
Type=simple
# Edit this path if you clone elsewhere.
WorkingDirectory=%h/git/jarvis-router
ExecStart=/usr/bin/env bash %h/git/jarvis-router/scripts/forge-watcher.sh
Restart=always
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=forge-watcher

[Install]
WantedBy=default.target
